# docker/cron/crontab/docker

1 0 * * *  /usr/local/bin/php /var/www/artisan import:cdfs >> /var/log/cron.log
* * * * *  /bin/echo "foobar" >> /var/log/cron.log

