# docker/cron/crontab/docker

1 0 * * *  /usr/local/bin/php /var/www/artisan import:web >> /var/log/cron.log
1 1 1 1 1  /bin/echo "foobar" >> /var/log/cron.log

